# root/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:aircast]
redirect_stderr=true
command=/bin/aircast-linux-x86_64
process_name = aircast-linux-x86_64
